# supervisor/book_seller_travis.conf
[program:book_seller]
command=`which gunicorn` book_seller.wsgi:application
directory=~/book_seller/
user=travis
autostart=true
autorestart=true
redirect_stderr=true