# supervisor/justbookcrawler.conf
[program:justbookcrawler]
command=/home/vagrant/.virtualenvs/book_seller/bin/scrapyd
directory=/vagrant/justbookcrawler
user=vagrant
autostart=true
autorestart=true
redirect_stderr=true