# supervisor/justbookcrawler_travis.conf
[program:justbookcrawler]
command=`which scrapyd`
directory=/travis/justbookcrawler
user=travis
autostart=true
autorestart=true
redirect_stderr=true